# templates/status.cron
#!/bin/bash

source /home/seqware/.bash_profile

seqware workflow-run launch-scheduled --host `hostname --long`
seqware workflow-run propagate-statuses --host `hostname --long` --threads 10

